# scripts/aether.service
[Unit]
Description=Aether Tool Collection
After=network.target

[Service]
Type=simple
User=novice
WorkingDirectory=/home/novice/project/Aether
ExecStart=/usr/bin/npm start -- -p 3001
Restart=on-failure
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
